# deploy/systemd/uo_regulations_translate_worker.service.template
[Unit]
Description=PDF OCR Translate Worker
After=network.target

[Service]
Type=simple
User={{APP_USER}}
Group=www-data
WorkingDirectory={{APP_ROOT}}
EnvironmentFile={{ENV_FILE}}
Environment="PATH={{APP_ROOT}}/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart={{APP_ROOT}}/.venv/bin/python worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
